# gonetmon.conf
start on startup

script
   /usr/local/bin/gonetmon  -device="enp4s6" -cidr="192.168.2.0/24" >> /var/log/gonetmon.log
end script

